# scripts/deploy/contra-signaling.service
[Unit]
Description=Contra WebRTC Signaling
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/contra/apps/signaling
Environment=PORT=8080
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
# 203/EXEC：node 路径不对；1/FAILURE：缺 dist 或 npm 依赖 — 见 journalctl
ExecStart=/opt/contra/scripts/deploy/start-signaling.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
