# scripts/deploy/contra-signaling.service
[Unit]
Description=Contra WebRTC Signaling
After=network.target

[Service]
Type=simple
# OpenCloudOS / RHEL 无 www-data；轻量 VPS 用 root 即可
User=root
Group=root
WorkingDirectory=/opt/contra/apps/signaling
Environment=PORT=8080
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
