# scripts/deploy/contra-signaling.service
[Unit]
Description=Contra WebRTC Signaling
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/contra/apps/signaling
Environment=PORT=8080
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
