# scripts/deploy/contra-signaling.service
[Unit]
Description=Contra WebRTC Signaling
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/contra/apps/signaling
Environment=PORT=8080
Environment=NODE_ENV=production
# 203/EXEC 常见原因：/usr/bin/node 不存在；用脚本自动查找 node
ExecStart=/opt/contra/scripts/deploy/start-signaling.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
