# scripts/deploy/contra-signaling.service
[Unit]
Description=Contra WebRTC Signaling
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/contra/apps/signaling
Environment=PORT=8080
Environment=NODE_ENV=production
Environment=HOME=/root
Environment=NVM_DIR=/root/.nvm
# nvm 用户：复制 signaling.env.example → /etc/contra/signaling.env 并设置 NODE_BIN
EnvironmentFile=-/etc/contra/signaling.env
StandardOutput=journal
StandardError=journal
ExecStart=/opt/contra/scripts/deploy/start-signaling.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
